# discord-bot.service
# Copy to /etc/systemd/system/discord-bot.service, then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now discord-bot
# /usr/local/bin/bun is a stable symlink to the installed bun (runs .ts natively
# and auto-loads .env from WorkingDirectory):
#   sudo ln -sf "$(which bun)" /usr/local/bin/bun   (rerun after bun upgrades)
[Unit]
Description=Discord bot
After=network-online.target
Wants=network-online.target

[Service]
User=ftrookie
WorkingDirectory=/home/ftrookie/discord-bot
ExecStart=/usr/local/bin/bun run src/Index.ts
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
